# 03｜CHEN操盘手系统/01｜CHEN外接大脑/作品链接采集器/chen-content-link-collector.service.example
[Unit]
Description=CHEN Content Link Collector Webhook
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/chen-content-link-collector
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
