# 语料库/00｜素材管理自动同步.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chenbrain.sync-materials-to-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/chen.zip/Desktop/CHEN BRAIN/语料库/00｜同步素材管理到日常语料.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/chen.zip/Desktop/CHEN BRAIN</string>

  <key>WatchPaths</key>
  <array>
    <string>/Users/chen.zip/Desktop/CHEN BRAIN/语料库/素材管理</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/chen.zip/Desktop/CHEN BRAIN/语料库/.素材管理自动同步.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/chen.zip/Desktop/CHEN BRAIN/语料库/.素材管理自动同步.err.log</string>
</dict>
</plist>
