# 03｜CHEN操盘手系统/03｜CHEN语料库/00｜素材管理自动同步.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chenbrain.sync-materials-to-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>CHEN_BRAIN_ROOT='/Users/chen.zip/Desktop/CHEN BRAIN' /usr/bin/python3 '/Users/chen.zip/Library/Application Support/CHEN-BRAIN/sync_materials_to_daily.py'; sleep 10</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CHEN_BRAIN_ROOT</key>
    <string>/Users/chen.zip/Desktop/CHEN BRAIN</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/chen.zip/Desktop/CHEN BRAIN</string>

  <key>WatchPaths</key>
  <array>
    <string>/Users/chen.zip/Desktop/CHEN BRAIN/语料库/素材管理</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/chen.zip/Desktop/CHEN BRAIN/语料库/.素材管理自动同步.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/chen.zip/Desktop/CHEN BRAIN/语料库/.素材管理自动同步.err.log</string>
</dict>
</plist>
