# 03｜CHEN操盘手系统/02｜CHEN语料库/00｜素材管理自动同步.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chenbrain.sync-materials-to-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/chen.zip/Library/Application Support/CHEN-BRAIN/watch_materials_to_daily.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/chen.zip/Library/Application Support/CHEN-BRAIN</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/chen.zip/Library/Application Support/CHEN-BRAIN/sync_materials_to_daily.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/chen.zip/Library/Application Support/CHEN-BRAIN/sync_materials_to_daily.err.log</string>
</dict>
</plist>
